# cmd/freegeoip/deploy/roles/redis/files/rc.local
#!/bin/sh -e

iptables-restore < /etc/iptables.conf

exit 0
