# cmd/freegeoip/ansible-playbook/roles/prometheus/templates/upstart.conf
description "prometheus web server"

start on (started networking)
stop on runlevel [!2345]

limit nofile 200000 200000

respawn limit unlimited
chdir /opt/prometheus
exec ./prometheus -config.file=./prometheus.yml {{ prometheus_params | join(" ") }}
post-stop exec sleep 1
