# cmd/freegeoip/ansible-playbook/roles/nodeexporter/files/upstart.conf
description "node exporter web server"

start on (started networking)
stop on runlevel [!2345]

respawn limit unlimited
exec /usr/local/bin/node_exporter
post-stop exec sleep 1
